# deploy/softportal-api.http-ip.service
[Unit]
Description=SoftPortal FastAPI backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/SoftPortal/backend
EnvironmentFile=/srv/SoftPortal/backend/.env
ExecStart=/srv/SoftPortal/backend/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=3
User=softportal
Group=softportal

[Install]
WantedBy=multi-user.target
